# crontab
0 8-18 * * 1-6 cd / && python3 /main.py >> /logfile 2>&1