# crontab
0 8-18 * * * cd / && python3 /main.py >> /logfile