# crontab
0 8-18 * * * python /main.py >> /logfile
